# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command = service nginx start 
user = root
autostart = true

[program:php]
command = service php7.0-fpm start
user = root
autostart = true